# python_exec/rootfs/etc/services.d/python_executor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the python executor service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare code, requirements

code=$(bashio::config 'code')
requirements=$(bashio::config 'requirements')

bashio::log.info "${code}"

bashio::log.info "${requirements}"

PYTHON=$(which python3)

virtualenv -p ${PYTHON} .

source bin/activate

if [ -n "$requirements" ]; then
  pip install -U ${requirements}
fi

python3 ${code}
