# python_exec/rootfs/etc/services.d/python_exec/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the python executor service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare code
declare requirements

code=$(bashio::config 'code')
requirements=$(bashio::config 'requirements')

PYTHON=$(which python3)

bashio::log.info "Create venv"
${PYTHON} -m venv .

bashio::log.info "Enable venv"
source bin/activate

if [ -n "$requirements" ]; then
  bashio::log.info "Install deps: ${requirements}"
  pip install -U ${requirements}
fi

python3 <<EOF
${code}
EOF
