# python_exec/rootfs/etc/services.d/python_exec/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the python executor service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

PYTHON=$(which python3)

bashio::log.info "Create venv"
${PYTHON} -m venv .

bashio::log.info "Enable venv"
source bin/activate

if bashio::config.has_value "requirements"; then
  bashio::log.info "Install requirements"
  yes | pip install --disable-pip-version-check -q -U $(bashio::config 'requirements')
fi

bashio::log.info "Run source code"
python3 <<EOF
$(bashio::config 'code')
EOF

bashio::exit.ok