# python_exec/rootfs/etc/services.d/python_exec/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the python executor service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

PYTHON=$(which python3)

if bashio::config.has_value "preexec"; then
  bashio::log.info "Run user pre-command..."
  bash -c "$(bashio::config 'preexec')"
fi

if bashio::config.has_value "requirements"; then
  bashio::log.info "Install Python requirements..."
  pip install -q -U pip
  pip install -q -U $(bashio::config 'requirements')
fi

bashio::log.info "Run source code..."
python3 <<EOF
$(bashio::config 'code')
EOF
